# profile.d/mail.sh
export MAIL=~/Maildir
export MAIL=/var/vmail/${USER}/Maildir
export MAIL=/var/vmail/$USER/Maildir
